# packages/binary-kde-settings_12.0.0_amd64/etc/profile.d/about-distro.sh
 
cp /usr/share/binary-kde-settings/kcm-about-distrorc /etc/xdg/
